# etc/cron.d/elb_update_cert
# min hour day month dow (0 Sunday)
34 20 * * * root bash /home/ec2-user/elb_update_cert.sh |& logger -t elb_update_cert -p local0.info
